# crontab
# m h  dom mon dow   command
* * * * * sh /var/www/html/fortuneCowsay.sh
~                    
